# services/php-fpm/run
#!/bin/bash

exec /usr/sbin/php-fpm${OC_PHP_VERS} -F -y /config/nginx/nginx-fpm.conf
