# services/mariadb/run
#!/bin/bash

# Function to handle clean shutdown
shutdown_mariadb() {
	if pgrep -x mariadbd >/dev/null 2>&1; then
		echo "Shutting down MariaDB cleanly..."
		mysqladmin --socket=/run/mysqld/mysqld.sock -u root shutdown >/dev/null 2>&1
	else
		echo "MariaDB is not running. Skipping shutdown."
	fi
}

# Trap SIGTERM and call shutdown function
trap 'shutdown_mariadb' SIGTERM

# Start MariaDB in the background
/usr/bin/mariadbd-safe --user=abc --skip-syslog --binlog-format=MIXED &

# Wait for the MariaDB process to finish (keeps the container running)
wait $!
